# debian/prerm
#!/bin/sh

set -e

case "$1" in
    remove|purge)
        invoke-rc.d solr5 stop
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

exit 0
